# supervisord/bot_de.ini
[program:kann_ai_bot_de]
; Nach ~/.config/supervisord/conf.d/kann_ai_bot_de.ini kopieren.

directory      = /var/www/virtual/%(ENV_USER)s/kann_ai_bot
command        = /var/www/virtual/%(ENV_USER)s/kann_ai_bot/.venv/bin/python bot.py
autostart      = true
autorestart    = true
startsecs      = 5
stopwaitsecs   = 30
environment    = DOTENV_PATH="/var/www/virtual/%(ENV_USER)s/kann_ai_bot/.env.de"
stdout_logfile          = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_bot_de.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stderr_logfile          = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_bot_de.err.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups  = 5
