# supervisord/bot_de.ini
[program:kann_ai_bot_de]
directory      = REPO_PATH
command        = REPO_PATH/.venv/bin/python bot.py
autostart      = true
autorestart    = true
startsecs      = 5
stopwaitsecs   = 30
environment    = DOTENV_PATH="REPO_PATH/.env.de"
stdout_logfile          = LOG_PATH/kann_ai_bot_de.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stderr_logfile          = LOG_PATH/kann_ai_bot_de.err.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups  = 5
