# supervisord/web.ini
[program:kann_ai_web]
; Nach ~/.config/supervisord/conf.d/kann_ai_web.ini kopieren.
; Vorher Port reservieren:
;   uberspace web backend set / --http --port PORT
; Dann PORT hier eintragen.

directory      = /var/www/virtual/%(ENV_USER)s/kann_ai_bot
command        = /var/www/virtual/%(ENV_USER)s/kann_ai_bot/.venv/bin/uvicorn
                 web.app:app --host 127.0.0.1 --port PORT
autostart      = true
autorestart    = true
startsecs      = 5
stopwaitsecs   = 10
stdout_logfile = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_web.log
stderr_logfile = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_web.err.log
