# supervisord/web_de.ini
[program:kann_ai_web_de]
; Nach ~/.config/supervisord/conf.d/kann_ai_web_de.ini kopieren.
; Vorher Port reservieren:
;   uberspace web backend set / --http --port PORT_DE
; Dann PORT_DE hier eintragen.

directory      = /var/www/virtual/%(ENV_USER)s/kann_ai_bot
command        = /var/www/virtual/%(ENV_USER)s/kann_ai_bot/.venv/bin/uvicorn
                 web.app:app --host 127.0.0.1 --port PORT_DE
autostart      = true
autorestart    = true
startsecs      = 5
stopwaitsecs   = 10
environment    = DOTENV_PATH="/var/www/virtual/%(ENV_USER)s/kann_ai_bot/.env.de"
stdout_logfile          = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_web_de.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 5
stderr_logfile          = /var/www/virtual/%(ENV_USER)s/logs/kann_ai_web_de.err.log
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups  = 5
